# install-scripts/6-desktop/files/firstRun.service
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/firstrun
User=root

[Unit]
Before=lightdm.service

[Install]
WantedBy=multi-user.target
